# copied-to-kiosk/opt/kiosk/template/etc/systemd/system/gdm.service.d/after_kiosk-updater.conf
[Unit]
Wants=kiosk-updater.service
After=kiosk-updater.service
